# root/etc/services.d/nut/run
#!/usr/bin/with-contenv bash

cd /app/nut || exit

exec \
	s6-setuidgid abc python3 nut.py -S